# root/etc/services.d/qbittorrent/run
#!/usr/bin/with-contenv bash

# 设置权限掩码umask

UMASK_SET=${UMASK_SET:-000}
umask $UMASK_SET

# 启动qBittorrent
exec \
	s6-setuidgid abc qbittorrent-nox --webui-port=$WEBUI_PORT --profile=/config
